# rootfs/etc/systemd/system/getty@tty1.service.d/zzz-sway-autologin.conf
# Autologin the desktop user on tty1. The login shell launches sway (see the
# sway user's ~/.bash_profile, installed by the Dockerfile). pam_systemd then
# registers a logind session on seat0, which grants DRM/input device access.
#
# Named to sort AFTER Kairos's live-mode "override.conf" (which autologins root
# on tty1), so our sway autologin wins. systemd merges getty drop-ins in
# alphabetical filename order; the last ExecStart= wins.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin sway --noclear %I $TERM
