# rootfs/usr/lib/systemd/system/hadron-user-setup.service
[Unit]
Description=Hadron per-user setup (flatpak Flathub remote)
# Run after Kairos has applied the cloud-config: /etc/passwd is not persisted, so
# the desktop user is recreated every boot by the cos-setup boot stage. Wait for
# cos-setup-boot (and -network) to finish first, else the user doesn't exist yet.
After=cos-setup-network.service cos-setup-boot.service systemd-user-sessions.service
ConditionPathExists=/etc/passwd

[Service]
Type=oneshot
ExecStart=/usr/bin/hadron-user-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
