# rootfs/usr/lib/systemd/system/hadron-splash.service
[Unit]
Description=Hadron boot splash
# hadron-desktop runs ly (not getty) on tty1, so we can't mirror the upstream
# splash unit's Before=getty.target. Instead order this oneshot Before=ly@tty1
# and have it pulled into multi-user.target alongside ly (see the Dockerfile
# wiring): systemd then runs the splash first and ly waits for it to finish,
# giving the VT a ~5s animated "HADRON" before the login screen takes over.
# Gated off the live/installer boot — there the installer owns tty1 — with the
# same conditions as ly's own ly@.service drop-in.
Before=ly@tty1.service
ConditionPathExists=/usr/bin/hadron-splash
ConditionKernelCommandLine=!install-mode
ConditionPathExists=!/run/cos/live_mode

[Service]
Type=oneshot
# Stay "active (exited)" after the splash finishes so a second `systemctl start`
# (the 91_sway_session boot stage also starts it, to guarantee splash-before-ly
# even if the baked multi-user.target.wants symlink is shadowed by the /etc
# overlay) is a no-op instead of replaying the animation.
RemainAfterExit=yes
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
ExecStart=/usr/bin/hadron-splash
# hadron-splash self-limits to ~5s and handles SIGTERM; bound the worst case so a
# hang can never stall ly (a oneshot blocks its Before= successors).
TimeoutStartSec=10s

[Install]
WantedBy=multi-user.target
