# rootfs/etc/systemd/system/ly.service.d/10-conditions.conf
# ly is wired into multi-user.target (Kairos forces `set-default multi-user.target`,
# so graphical.target is never reached). But multi-user.target is also reached
# during the live/installer boot — and ly would grab tty1, which the installer
# (kairos-installer.service) needs. Gate ly so it only starts on a normally
# booted, installed system: not during install-mode, and not on the live ISO.
[Unit]
ConditionKernelCommandLine=!install-mode
ConditionPathExists=!/run/cos/live_mode
